# NVIDIA-Linux-x86_64-460.73.01-grid/init-scripts/upstart/nvidia-gridd.conf
# NVIDIA Grid Daemon Init Script
#
# Copyright (c) 2015 NVIDIA Corporation
#
# All rights reserved.  All information contained herein is proprietary and
# confidential to NVIDIA Corporation.  Any use, reproduction, or disclosure
# without the written permission of NVIDIA Corporation is prohibited.
#

description "NVIDIA Grid Daemon"

start on runlevel [2345]
stop on runlevel [016]

env NVGD_BIN=/usr/bin/nvidia-gridd
env NVGD_RUNTIME=/var/run/nvidia-gridd
env LD_LIBRARY_PATH=/usr/lib/nvidia/gridd
expect fork

exec $NVGD_BIN

post-stop script
    rm -rf $NVGD_RUNTIME
end script
